# distrib/initscripts/netbsd.macipgw.in
#!/bin/sh
#
# PROVIDE: macipgw
# REQUIRE: atalkd
#
# IP-over-AppleTalk (netatalk @netatalk_version@)
# Adjust the arguments to match your network.

. /etc/rc.subr

name="macipgw"
rcvar=$name
command="@sbindir@/macipgw"
command_args="-n 8.8.8.8 192.168.151.0 255.255.255.0"

load_rc_config $name
run_rc_command "$1"
