# distrib/initscripts/netbsd.netatalk.in
#!/bin/sh
#
# PROVIDE: netatalk
# KEYWORD: shutdown
#
# Netatalk AFP fileserver (netatalk @netatalk_version@)

. /etc/rc.subr

name="netatalk"
rcvar=$name
command="@sbindir@/netatalk"
etcdir="@etcdir@"
pidfile="/var/run/${name}.pid"
required_files="$etcdir/afp.conf"

load_rc_config $name
run_rc_command "$1"
