# distrib/initscripts/netbsd.a2boot.in
#!/bin/sh
#
# PROVIDE: a2boot
# REQUIRE: atalkd
#
# AppleTalk Apple II netboot daemon (netatalk @netatalk_version@)

. /etc/rc.subr

name="a2boot"
rcvar=$name
command="@sbindir@/a2boot"

load_rc_config $name
run_rc_command "$1"
