# distrib/initscripts/netbsd.papd.in
#!/bin/sh
#
# PROVIDE: papd
# REQUIRE: atalkd
#
# AppleTalk print server daemon (netatalk @netatalk_version@)

. /etc/rc.subr

name="papd"
rcvar=$name
command="@sbindir@/papd"
pidfile="/var/run/${name}.pid"
required_files="@etcdir@/papd.conf"

load_rc_config $name
run_rc_command "$1"
