# distrib/initscripts/netbsd.timelord.in
#!/bin/sh
#
# PROVIDE: timelord
# REQUIRE: atalkd
#
# AppleTalk time server daemon (netatalk @netatalk_version@)

. /etc/rc.subr

name="timelord"
rcvar=$name
command="@sbindir@/timelord"

load_rc_config $name
run_rc_command "$1"
