# distrib/initscripts/openbsd.netatalk.in
#!/bin/ksh
#
# Netatalk @netatalk_version@ daemons.
#

daemon="@sbindir@/netatalk"

. /etc/rc.d/rc.subr

rc_start() {
        rc_exec ". ~/.profile; ${daemon}"
}

rc_cmd $1
