# distrib/initscripts/openrc.papd.in
#!/sbin/openrc-run

# Netatalk @netatalk_version@ init script.
# PAP print server (papd)

name=$RC_SVCNAME
command="@sbindir@/${RC_SVCNAME}"

depend() {
		need net
		need atalkd
}
