# distrib/initscripts/rc.netbsd.tmpl
#!/bin/sh
#
# PROVIDE: netatalk
# KEYWORD: shutdown
#
# AFP fileserver for Macintosh clients (netatalk @netatalk_version@)
#

. /etc/rc.subr

name="netatalk"
rcvar=$name
command="@sbindir@/netatalk"
etcdir="@etcdir@"
pidfile="/var/run/${name}.pid"
required_files="$etcdir/afp.conf"

load_rc_config $name
run_rc_command "$1"
