# distrib/initscripts/systemd.papd.service.in
# This file is part of Netatalk @netatalk_version@.

[Unit]
Description=AppleTalk PAP printer server daemon
Documentation=man:papd(8) man:papd.conf(5)
Documentation=https://netatalk.io/
After=network-online.target atalkd.service
Requires=atalkd.service

[Service]
Type=forking
GuessMainPID=no
ExecStart=@sbindir@/papd
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
