# distrib/debian/netatalk.postrm
#!/bin/sh
set -e

if [ $1 = "purge" ]; then
  update-rc.d netatalk remove > /dev/null
# update-rc.d afpd remove > /dev/null
# update-rc.d papd remove > /dev/null
fi
